# packaging/rpm/systemprompt.service
[Unit]
Description=systemprompt AI governance gateway
Documentation=https://systemprompt.io/documentation/
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=systemprompt
Group=systemprompt
EnvironmentFile=-/etc/systemprompt/systemprompt.env
ExecStart=/usr/bin/systemprompt
Restart=on-failure
RestartSec=5s
LimitNOFILE=65535

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/systemprompt /var/log/systemprompt

[Install]
WantedBy=multi-user.target
